# infra/supervisord.ini
[supervisord]
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=100MB
logfile_backups=10
loglevel=info
nodaemon=true

[supervisorctl]
serverurl=http://localhost:9001

[program:poca-main]
command=run-program python3.12 -m src --host 0.0.0.0 --port 8000
directory=/
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
stdout_logfile=/var/log/poca/main.stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stdout_syslog=true
stderr_logfile=/var/log/poca/main.stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stderr_syslog=true
environment=REDIS_DSN="%(ENV_REDIS_DSN)s"

[program:poca-session-cleaner]
command=run-program python3.12 -m src.cli session-cleaner
directory=/
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=INT
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/poca/session-cleaner.stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stdout_syslog=true
stderr_logfile=/var/log/poca/session-cleaner.stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stderr_syslog=true
environment=REDIS_DSN="%(ENV_REDIS_DSN)s"

[program:poca-scanner-manager]
command=run-program python3.12 -m src.cli scanner-manager
directory=/
autostart=true
autorestart=true
startsecs=10
startretries=5
stopsignal=INT
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=root
stdout_logfile=/var/log/poca/scanner-manager.stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stdout_syslog=true
stderr_logfile=/var/log/poca/scanner-manager.stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stderr_syslog=true
environment=REDIS_DSN="%(ENV_REDIS_DSN)s"
